# infra/hetzner/templates/{{PROJECT_NAME}}-<job>.timer.j2
[Unit]
Description={{PROJECT_NAME}} <job> schedule

[Timer]
OnCalendar=hourly
Persistent=true
Unit={{PROJECT_NAME}}-<job>.service

[Install]
WantedBy=timers.target
